# automation/launchd/com.jarvis.tg-media-analyzer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.tg-media-analyzer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/seregaateist/Projects/jarvis/modules/tg-media-analyzer/venv/bin/python3.11</string>
        <string>/Users/seregaateist/Projects/jarvis/modules/tg-media-analyzer/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seregaateist/Projects/jarvis/modules/tg-media-analyzer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/seregaateist/Projects/jarvis/modules/tg-media-analyzer/logs/bot.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seregaateist/Projects/jarvis/modules/tg-media-analyzer/logs/bot.stdout.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
