# automation/launchd/com.jarvis.rafail-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.jarvis.rafail-cron</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/seregaateist/Projects/jarvis/modules/rafail/scripts/run_collect.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/seregaateist/Projects/jarvis</string>
    <key>StartInterval</key><integer>21600</integer>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key>
    <string>/Users/seregaateist/Projects/jarvis/logs/rafail-cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seregaateist/Projects/jarvis/logs/rafail-cron.log</string>
</dict>
</plist>
