# automation/launchd/com.jarvis.process-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.process-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/Users/seregaateist/Projects/jarvis/infra/process_manager.py</string>
        <string>watch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seregaateist/Projects/jarvis</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/seregaateist/Projects/jarvis/logs/process_manager.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seregaateist/Projects/jarvis/logs/process_manager.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>JARVIS_ROOT</key>
        <string>/Users/seregaateist/Projects/jarvis</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
